# deploy/com.colors-agent.colors.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.colors-agent.colors</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/YOUR_USERNAME/.colors-agent/dist/cli.js</string>
    <string>channel</string>
    <string>telegram</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>ANTHROPIC_API_KEY</key>
    <string>YOUR_KEY_HERE</string>
    <key>COLORS_PASSPHRASE</key>
    <string>YOUR_PASSPHRASE_HERE</string>
    <key>TELEGRAM_BOT_TOKEN</key>
    <string>YOUR_TOKEN_HERE</string>
    <key>TELEGRAM_ALLOWED_USERS</key>
    <string>*</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/colors-agent.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/colors-agent-error.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
